# crontab
0 0 * * * cd /Downloads/ProjetoLAPR_ESOF/ && php artisan schedule:run >> /dev/null 2>&1
